# deploy/systemd/fiochat-telegram.service
[Unit]
Description=Fiochat Telegram Bot Bridge
After=network-online.target fiochat.service
Wants=network-online.target fiochat.service

[Service]
Type=simple
User=svc
WorkingDirectory=/opt/fiochat/telegram
Environment=NODE_ENV=production

# Keep config consistent across services
Environment="FIOCHAT_CONFIG_FILE=/etc/fiochat/config.yaml"
Environment="HOME=/var/lib/fiochat"

# Load secrets / bot token etc without editing the unit
EnvironmentFile=-/etc/fiochat/telegram.env

ExecStart=/usr/bin/node dist/index.js
Restart=on-failure
RestartSec=5

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only

StateDirectory=fiochat
ReadWritePaths=/var/lib/fiochat

[Install]
WantedBy=multi-user.target
