# deploy/systemd/fiochat.service
[Unit]
Description=Fiochat AI Service (HTTP Server)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=svc
WorkingDirectory=/opt/fiochat

# Prefer system-wide config for production
Environment="AICHAT_CONFIG=/etc/fiochat/config.yaml"

# Deterministic HOME that is writable even with ProtectHome=read-only
Environment="HOME=/var/lib/fiochat"

ExecStart=/usr/local/bin/fiochat --serve 127.0.0.1:8000
Restart=on-failure
RestartSec=5

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only

# systemd-managed writable directories (created automatically)
StateDirectory=fiochat

# Only allow writes where we expect them
ReadWritePaths=/var/lib/fiochat

[Install]
WantedBy=multi-user.target
