# telegram/fio-telegram.service
[Unit]
Description=FIO Telegram server bot
After=network-online.target
# Add your AI service unit here if needed:
# After=network-online.target <your-ai-service>.service
# Wants=<your-ai-service>.service

[Service]
Type=simple
WorkingDirectory=/opt/fio-bot
Environment=NODE_ENV=production
ExecStart=/usr/bin/npm run start
Restart=on-failure
RestartSec=5
User=svc

[Install]
WantedBy=multi-user.target